# local.mwm.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.mwm</string>

  <key>ProgramArguments</key>
  <array>
    <string>@MWM_BIN@</string>
    <string>daemon</string>
    <string>--columns</string>
    <string>2.5</string>
  </array>

  <key>WorkingDirectory</key>
  <string>@MWM_WORKDIR@</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/mwm.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/mwm.err.log</string>
</dict>
</plist>
